# deploy/champion-sglang.service
[Unit]
Description=Champion SGLang inference server (Docker)
After=network-online.target docker.service
Requires=docker.service
Wants=network-online.target

[Service]
# SGLang runs in the lmsysorg/sglang image: its bundled cuda-compat lets the
# CUDA-13 runtime work on the host's 535 driver (forward-compat on the L40).
# Serves the pre-merged champion full weights at /opt/champion on loopback :30000.
ExecStartPre=-/usr/bin/docker rm -f champion-sglang
ExecStart=/usr/bin/docker run --rm --name champion-sglang --gpus all \
  -p 127.0.0.1:30000:30000 \
  -v /ephemeral/champion:/opt/champion:ro \
  -v /ephemeral/hf_cache:/root/.cache/huggingface \
  lmsysorg/sglang:latest \
  python3 -m sglang.launch_server \
    --model-path /opt/champion --served-model-name champion \
    --host 0.0.0.0 --port 30000 \
    --tensor-parallel-size 1 --dtype float16 \
    --enable-deterministic-inference --random-seed 0 \
    --tool-call-parser qwen25
ExecStop=/usr/bin/docker stop champion-sglang
Restart=always
RestartSec=10
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
